# fabric/deploy/supervisor.conf
[group:%(proj_name)s]
programs=nginx

[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true
redirect_stderr=true
